# host/supervisord.devcontainer.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid

[program:websock]
command=python3 -u /app/pythings/websock_server/ws_server.py
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=10

# Dev-only stand-in for the host-native TCP-to-serial bridge. There is no
# dongle (and nobody starts the real host-side bridge) in the dev container,
# so this fakes a healthy dongle on TCP :9000 -- otherwise firework-daemon
# spams "Connection refused" forever and the UI pins the dongle at
# DEVICE_ERROR. Never present in supervisord.conf (production talks to the
# real bridge). Starts before firework-daemon (priority 15 < 20) so the
# daemon's first connect lands. Stop it (supervisorctl stop mock-bridge) and
# set BYH_BRIDGE_HOST back to host.docker.internal to use real hardware.
[program:mock-bridge]
command=python3 -u /app/pythings/mock_bridge/mock_bridge.py
environment=BYH_BRIDGE_BIND=127.0.0.1,BYH_BRIDGE_PORT=9000
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=15

[program:firework-daemon]
command=python3 -u /app/pythings/pc_daemon/pc_daemon.py
# Dial the in-container mock bridge (above) instead of host.docker.internal.
environment=BYH_BRIDGE_HOST=127.0.0.1,BYH_BRIDGE_PORT=9000
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=20

[program:audio-player]
command=python3 -u /app/pythings/audio_player/audio_player.py
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=25

[program:byh-app]
command=npm run dev
environment=PORT=1776,NODE_ENV=development,HOSTNAME=0.0.0.0
directory=/app/byh_app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=30